# com.screentime-check.collect.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.screentime-check.collect</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/robinerd/Documents/screentime-check/collect.py</string>
    </array>

    <!-- Run every hour -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- Also run on wake from sleep -->
    <key>StartOnMount</key>
    <false/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Users/robinerd/Documents/screentime-check/collect.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/robinerd/Documents/screentime-check/collect.log</string>
</dict>
</plist>
